# timeedit.service
[Unit]
Description=TimeEdit bot

[Service]
ExecStart=/opt/timeedit/timeedit.py

[Install]
WantedBy=default.target
